# files/supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=/usr/sbin/squid -NsY -f /etc/squid/squid.conf
stdout_events_enabled=true
stderr_events_enabled=true
